# lighting/cowrie/honeyfs/etc/systemd/system/lights.service
[Unit]
Description=Lighting Control Service

[Service]
ExecStart=/etc/lighting/run.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target