# appliances/cowrie/honeyfs/etc/systemd/system/appliances.service
[Unit]
Description=Appliance Control Service

[Service]
ExecStart=/etc/appliances/run.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
