# thermostat/cowrie/honeyfs/etc/systemd/system/thermostat.service
[Unit]
Description=Thermostat Control Service

[Service]
ExecStart=/etc/thermostat/run.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
